# webhook-listener.service
[Unit]
Description=Webhook Listener - Infrastructure Predictor CI/CD
After=network.target

[Service]
Type=simple
User=Adrian
WorkingDirectory=C:/Users/Adrian/infrastructure-predictor

# Simple Python HTTP server that listens for webhook POSTs
# and triggers deploy.sh. Can be replaced with GitHub webhooks, Gitea, etc.
ExecStart=C:/Users/Adrian/AppData/Local/Programs/Python/Python312/python.exe -X utf8 -c "
import http.server
import subprocess
import sys
import os

PORT = 9000
DEPLOY_SCRIPT = 'deploy.sh'

class WebhookHandler(http.server.BaseHTTPRequestHandler):
    def do_POST(self):
        content_len = int(self.headers.get('Content-Length', 0))
        body = self.rfile.read(content_len) if content_len else b''

        print(f'[webhook] Received POST {self.path} ({content_len} bytes)')

        # Verify secret token (optional)
        secret = self.headers.get('X-Webhook-Secret', '')
        expected = os.environ.get('WEBHOOK_SECRET', '')
        if expected and secret != expected:
            self.send_response(403)
            self.end_headers()
            self.wfile.write(b'Forbidden: invalid secret')
            return

        # Trigger deploy
        try:
            result = subprocess.run(
                ['bash', DEPLOY_SCRIPT],
                capture_output=True, text=True, timeout=300,
                cwd=os.path.dirname(os.path.abspath(__file__))
            )
            self.send_response(200 if result.returncode == 0 else 500)
            self.end_headers()
            self.wfile.write(result.stdout.encode() + result.stderr.encode())
        except subprocess.TimeoutExpired:
            self.send_response(504)
            self.end_headers()
            self.wfile.write(b'Deploy timed out')

    def do_GET(self):
        self.send_response(200)
        self.end_headers()
        self.wfile.write(b'Webhook listener active. POST to trigger deploy.')

    def log_message(self, format, *args):
        print(f'[webhook] {args[0]}')

server = http.server.HTTPServer(('0.0.0.0', PORT), WebhookHandler)
print(f'Webhook listener started on port {PORT}')
server.serve_forever()
"

Restart=always
RestartSec=10

# Security: set a WEBHOOK_SECRET env var
Environment="WEBHOOK_SECRET=change-me-to-a-random-string"

[Install]
WantedBy=multi-user.target
